# init/data/alios/service/agentwatch.rhel5
#!/bin/bash
#
# agentwatch      xen/kvm gshell
#

# chkconfig: 2345 98 01
# description: Starts and stops guest agent

SNAME=agentwatch
PROG=/usr/sbin/gshelld
VIRTWHAT=/usr/sbin/virt-what-cpuid-helper
RESTART1=1
SHUT=0
which=""

shut() {
    if [ $SHUT -eq 1 ] || [ $SHUT -eq 0 ]; then
        for gshelld_pid in `ps axf | grep "$PROG" | grep -v grep | awk {'print $1'}`
        do
            kill -9 $gshelld_pid
        done
    fi
    if [ $SHUT -eq 2 ] || [ $SHUT -eq 0 ]; then
	result=`cat /etc/inittab |grep "aliyun-service"`
	if [ $? == 0 ];then
	    sed -i /aliyun-service/d /etc/inittab	    
	    kill -HUP 1
	    sleep 1
	fi
    fi
}

which()
{
	skip_qemu_kvm=false
	cpuid=$($VIRTWHAT)
	if [ "$cpuid" = "XenVMMXenVMM" ]; then 
	    skip_qemu_kvm=true
	fi
	if ! "$skip_qemu_kvm"; then
	    if [ "$cpuid" = "KVMKVMKVM" ]; then
		which="kvm"
	    fi
	fi
}

# start function
start() {
    sleep 5
    if [ -f /var/lock/$SNAME ]
    then
	result=`ps aux |grep -E 'aliyun-service|gshelld'`
	echo $result |grep '/usr/sbin' 1>/dev/null 2>&1
	if [ $? -eq 1 ]; then
    	    rm -rf /var/lock/$SNAME
	else
	    RESTART1=0
	    echo "exist"
	    echo "$SNAME is already started!" 1>/dev/null 2>&1
	fi
    fi

    which 

    result=`ps aux |grep -E 'aliyun-service|gshelld'`
    pid11=`echo $result |grep '/usr/sbin/aliyun-service' |awk '{print $2}'` 1>/dev/null 2>&1
    aliyun=$?
    pid21=`echo $result |grep '/usr/sbin/gshelld' |awk '{print $2}'` 1>/dev/null 2>&1
    gshelld=$?
    if [ $pid11 ] || [ $pid21 ];then
    sleep 1
    fi
    if [ $pid11 ]; then
        result=`ps aux |grep -E 'aliyun-service|gshelld'`
        pid12=`echo $result |grep '/usr/sbin/aliyun-service' |awk '{print $2}'` 1>/dev/null 2>&1
	if [ "$pid11" = "$pid12" ]; then
	    echo "the same stop gshelld"	
	    RESTART1=0
	    SHUT=1
	fi
    fi
    if [ $pid21 ]; then
        result=`ps aux |grep -E 'aliyun-service|gshelld'`
        pid22=`echo $result |grep '/usr/sbin/gshelld' |awk '{print $2}'` 1>/dev/null 2>&1
	if [ "$pid21" = "$pid22" ]; then
	    echo "the same stop aliyun"	1>/dev/null 2>&1
	    RESTART1=0
	    SHUT=2
	fi
    fi

    shut

    if [ $RESTART1 -eq 1 ]; then
        if [ "$which" == "kvm" ];then
	    echo "7:2345:respawn:/usr/sbin/aliyun-service" >> /etc/inittab
	    kill -HUP 1
        else
            if [ ! -f /proc/xen/xenbus ];then
                mount none -t xenfs /proc/xen
            fi
            $PROG
            [ $? -eq 0 ] && touch /var/lock/subsys/gshelld
        fi
	touch /var/lock/$SNAME
        echo "finished" 1>/dev/null 2>&1
        exit 0
    fi
}

#stop function
stop() {
    echo "Stopping $SNAME ..."
    for gshelld_pid in `ps axf | grep "$PROG" | grep -v grep | awk {'print $1'}`
    do
        kill -9 $gshelld_pid
    done
    if [ -f /var/lock/subsys/gshelld ];then
        rm -rf /var/lock/subsys/gshelld
    fi
    rm -rf /var/lock/$SNAME
}

case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
reload|restart)
  stop
  start
  ;;
status)
  status $SNAME
  ;;
*)
  echo $"Usage: $0 {start|stop|restart|status}"
  exit 1
esac
