# init/data/centos/service/agentwatch.rhel5
#!/bin/bash
#
# agentwatch      xen/kvm gshell
#

# chkconfig: 2345 98 01
# description: Starts and stops guest agent

SNAME=agentwatch
PROG=/usr/sbin/gshelld
VIRTWHAT=/usr/sbin/virt-what-cpuid-helper
GSHELLLOCK=/var/lock/subsys/gshelld
aliyun_service_name='aliyun-service'
aliyun_service_bin="/usr/sbin/$aliyun_service_name"
aliyun_service_opts='-d'
agent_cmd="$aliyun_service_bin $aliyun_service_opts"
which=""

shut() {   
    if [ -f $GSHELLLOCK ]; then
        echo "kill gshell" 1>/dev/null 2>&1

        for gshelld_pid in `ps axf | grep "$PROG" | grep -v grep | awk {'print $1'}`
        do
            kill -9 $gshelld_pid
        done
        rm -rf $GSHELLLOCK    
    fi
}

which()
{
	skip_qemu_kvm=false
	cpuid=$($VIRTWHAT)
	if [ "$cpuid" = "XenVMMXenVMM" ]; then 
	    skip_qemu_kvm=true
	fi
	if ! "$skip_qemu_kvm"; then
	    if [ "$cpuid" = "KVMKVMKVM" ]; then
		which="kvm"
	    fi
	fi
}

# start function
start() {
    sleep 5
    if [ -f /var/lock/$SNAME ]
    then
	result=`ps aux |grep 'aliyun-service'`
	echo $result |grep '/usr/sbin' 1>/dev/null 2>&1
	if [ $? -eq 1 ]; then
    	    rm -rf /var/lock/$SNAME
	else
	    echo "exist"
	    echo "$SNAME is already started!" 1>/dev/null 2>&1
	fi
    fi

    which 
    shut

    if [ "$which" == "kvm" ];then
	    echo "7:2345:respawn:/usr/sbin/aliyun-service" >> /etc/inittab
	    kill -HUP 1
    else
        if [ ! -f /proc/xen/xenbus ];then
            mount none -t xenfs /proc/xen
        fi
        $agent_cmd
    fi
	touch /var/lock/$SNAME
    echo "finished" 1>/dev/null 2>&1
    exit 0
}

#stop function
stop() {
    shut
    echo "Stopping $SNAME ..."
    result=`cat /etc/inittab |grep "aliyun-service"`
	if [ $? == 0 ];then
	    sed -i /aliyun-service/d /etc/inittab	    
	    kill -HUP 1
	    sleep 1
	fi
    rm -rf /var/lock/$SNAME
}

case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
reload|restart)
  stop
  start
  ;;
status)
  status $SNAME
  ;;
*)
  echo $"Usage: $0 {start|stop|restart|status}"
  exit 1
esac
