# init/data/centos/service/agentwatch.cent7
#!/bin/bash
#
# agentwatch      xen/kvm gshell
#

# chkconfig: 2345 98 01
# description: Starts and stops guest agent

SNAME=agentwatch
PROG=/usr/sbin/gshelld
VIRTWHAT=/usr/sbin/virt-what-cpuid-helper
FILE=/etc/systemd/system/aliyun.service
RESTART1=1
GSHELLLOCK=/var/lock/subsys/gshelld
aliyun_service_name='aliyun-service'
aliyun_service_bin="/usr/sbin/$aliyun_service_name"
aliyun_service_opts='-d'
agent_cmd="$aliyun_service_bin $aliyun_service_opts"
which=""

shut() {    
    if [ -f $GSHELLLOCK ]; then
        echo "kill gshell" 1>/dev/null 2>&1

        for gshelld_pid in `ps axf | grep "$PROG" | grep -v grep | awk {'print $1'}`
        do
            kill -9 $gshelld_pid
        done
        rm -rf $GSHELLLOCK    
    fi
}

which() {
    skip_qemu_kvm=false
    cpuid="KVMKVMKVM"
    if [ -f VIRTWHAT ]; then
      cpuid=$($VIRTWHAT)
    fi
    if [ "$cpuid" = "XenVMMXenVMM" ]; then 
        skip_qemu_kvm=true
    fi
    if ! "$skip_qemu_kvm"; then
        if [ "$cpuid" = "KVMKVMKVM" ]; then
    	which="kvm"
        fi
    fi
}

# start function
start() {
    sleep 5
    if [ -f /var/lock/$SNAME ]
    then
	result=`ps aux |grep 'aliyun-service'`
	echo $result |grep '/usr/sbin' 1>/dev/null 2>&1
	if [ $? -eq 1 ]; then
    	    rm -rf /var/lock/$SNAME
	else
	    RESTART1=0
	    echo "exist"
        echo "$SNAME is already started!" 1>/dev/null 2>&1
	fi
    fi

    which
    stop

    if [ $RESTART1 -eq 1 ]; then
		echo "start to restart new" 1>/dev/null 2>&1
        if [ "$which" == "kvm" ];then
	    	echo "kvm" 1>/dev/null 2>&1
            touch $FILE
            echo "[Unit]" >> $FILE
	    	echo "Description=auto run aliyunservice or agent" >> $FILE
	    	echo "[Service]" >> $FILE
	    	echo "Type=simple" >> $FILE
	    	echo "ExecStart=/usr/sbin/aliyun-service" >> $FILE
	    	echo "Restart=always " >> $FILE
        echo "KillMode=process " >> $FILE
	    	echo "RestartSec=1" >> $FILE
	    	echo "[Install]" >> $FILE
	    	echo "WantedBy=multi-user.target" >> $FILE
	    	/bin/systemctl daemon-reload  
            /bin/systemctl start aliyun.service
        else
	    	if [ ! -f /proc/xen/xenbus ];then
				mount none -t xenfs /proc/xen
	    	fi
	    	echo "in xen" 1>/dev/null 2>&1
            $agent_cmd
        fi
        touch /var/lock/$SNAME
        echo "finished"
        exit 0
    fi
}

#stop function
stop() {
    shut
    echo "Stopping $SNAME ..." 1>/dev/null 2>&1  
    if [ -f $FILE ];then
        result=`/bin/systemctl list-units |grep aliyun`
        if [ "$result"x != ""x ];then
            echo "Start to stop" 1>/dev/null 2>&1
            /bin/systemctl stop aliyun.service
	    	/bin/systemctl daemon-reload  
        fi
        rm $FILE
    fi
    rm -rf /var/lock/$SNAME
}

case "$1" in
start)
  start
  ;;
stop)
  stop
  ;;
reload|restart)
  stop
  start
  ;;
status)
  status $SNAME
  ;;
*)
  echo $"Usage: $0 {start|stop|restart|status}"
  exit 1
esac
